# kokorodoki-monitor.service
[Unit]
Description=Monitor kokorodoki.service memory and activity

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/kokorodoki_monitor.sh
